# launchd/com.denderson.naia-telegram-bot.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.{{DONO_SLUG}}.naia-telegram-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/naiarodrigues/.pyenv/versions/3.12.8/bin/python3</string>
        <string>/Users/naiarodrigues/naia-bot/bot.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/naiarodrigues/naia-bot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/naiarodrigues/naia-bot/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/naiarodrigues/naia-bot/logs/launchd-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
